# slimm.service
[Unit]
Description=SLiMM — stateless DRM/KMS login bootstrapper (SLiM lineage)
Documentation=file:/usr/share/doc/slimm/README.md
DefaultDependencies=no
After=systemd-user-sessions.service
Before=graphical.target
Conflicts=getty@tty1.service
# DefaultDependencies=no drops the implicit shutdown ordering, so on poweroff/
# reboot the compositor's exit would let Restart=always pop a fresh greeter
# mid-shutdown. Re-add the conflict + ordering so slimm is stopped (not
# restarted) the moment shutdown begins.
Conflicts=shutdown.target
Before=shutdown.target
# Crash-loop guard: if the greeter (or session) dies >5×/30s, give up and stop
# so the machine falls back to a usable tty instead of flickering forever.
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
ExecStart=/usr/bin/slimm
Type=simple
# slimm has no reaper: the greeter authenticates then exec()s the compositor in
# place (the main PID *becomes* the compositor). systemd is the supervisor —
# when the compositor exits (logout), Restart=always starts a fresh greeter.
# So no lingering slimm process during the session, and the default
# control-group KillMode correctly tears the session down on stop/restart.
StandardInput=tty
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=1
# Escape (UI_QUIT) exits the greeter with 42 — do NOT respawn it, so the user
# drops to a console login and stays there (logout/crash still restart via
# Restart=always). Without this, Restart=always re-grabs tty1 a second after
# Escape and fights the console login.
RestartPreventExitStatus=42
TTYPath=/dev/tty1
TTYReset=yes
# Compositor setsid()s + ignores SIGHUP; keep VT hangups from disturbing it.
TTYVHangup=no
UtmpIdentifier=tty1

[Install]
WantedBy=graphical.target
